# master-crontab.cron
#!/bin/sh

# Copyright (C) 2013 Jonathan Geisler

if [ `ls /home/gharwood/active-contests | wc -l` -gt 0 ] ; then
	for file in /home/gharwood/active-contests/*
	do
		dir=`basename $file`
		cd "/home/gharwood/public_html/$dir/judge"
		./cronScript.php > ./errorLog.txt &
	done
fi
